# root/etc/services.d/tinyproxy/run
#!/usr/bin/with-contenv bash
set -e

exec /usr/bin/tinyproxy -c /etc/tinyproxy/tinyproxy.conf -d 2>&1
