# root/etc/services.d/tinyproxy/run
#!/usr/bin/with-contenv bash
set -e

exec /usr/sbin/tinyproxy -d 2>&1
